# deploy/native/mcp-poc.service
[Unit]
Description=imcp2 MCP server (mcp-poc)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/imcp2
Environment=PORT=8000
Environment=PUBLIC_URL=__PUBLIC_URL__
Environment=RUST_LOG=info
ExecStart=/opt/imcp2/mcp-poc
Restart=always
RestartSec=3
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
